# files/etc/oray.sh
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

curl http://teasiu:Teasiu123@ddns.oray.com/ph/update?hostname=ecoo.5166.info&myip=ipaddress 2>&1


